# deploy/systemd/aqsp-vibe-research.target
[Unit]
Description=AQSP Vibe-Research persistent service group
Requires=aqsp-vibe-research-api.service aqsp-vibe-research-preview.service
After=aqsp-vibe-research-api.service aqsp-vibe-research-preview.service

[Install]
WantedBy=multi-user.target
